# crontab
#* * * * * root echo "Hello world" >> /var/log/cron.log 2>&1
*/1 * * * * root mongo localhost:3001/meteor /mongo-job.js 2>&1

